# supervisor-mongod.conf
[program:mongod]
command=/usr/bin/mongod -f /etc/mongod.conf 
environment=LC_ALL=C
user=mongodb
autostart=false
autorestart=unexpected
