# supervisor-tater.conf
[program:tater]
command=/bin/bash /tater.sh
user=root
autostart=true
autorestart=unexpected
