# supervisor-mongod.conf
[program:mongod]
command=/usr/bin/mongod -f /etc/mongodb.conf 
environment=LC_ALL=C
user=mongodb
autostart=true
autorestart=unexpected
